# etc/profile.d/dot-profile-setup_2.0_all.sh
. .profile_00000
#export HOME_ECssm=$(s.unified_setup-cfg -base)
#HOME_ECssm=${HOME_ECssm%/dot-profile-setup*}
#s.which local_shortcuts.dot 2>/dev/null 1>/dev/null && . local_shortcuts.dot
#. .profile_zzzzz
